# deploy/whisper-wrap.service
[Unit]
Description=whisper-wrap v2 (in-process faster-whisper + Gemini /ask)
Documentation=https://github.com/shdennlin/whisper-wrap
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Adjust to your deployment. The user/group SHALL exist and own the install dir.
User=whisper
Group=whisper
WorkingDirectory=/opt/whisper-wrap

# uv is responsible for the venv. Set PATH explicitly so the unit doesn't depend
# on the invoking shell's profile.
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/whisper-wrap/.env

ExecStart=/usr/local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 8000

# v2 drops the in-app auto-restart machinery; recovery is the supervisor's job.
Restart=on-failure
RestartSec=5

# Hardening (loosen as needed for your environment).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# faster-whisper + CT2 can take ~10-30s to load the model at startup; uvicorn
# does not signal "ready" until the lifespan completes.
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
